# sys-inspector-journal.service
[Service]
ExecStart=/usr/local/share/sys-inspector/src/collectors/journald-collector.sh
User=root