# sys-inspector-network.service
[Unit]
Description=Network collector

[Service]
ExecStart=/usr/local/share/sys-inspector/src/collectors/network-collector.sh
User=root