# sys-inspector-process.timer
[Unit]
Description=Collect process tree every 5 minutes

[Timer]
OnCalendar=*:0/5

[Install]
WantedBy=timers.target