# systemd/sys-inspector-manifest.timer
[Unit]
Description=Sys-Inspector Service Manifest Weekly Timer

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
